# rulesetc/rules.d/70-snap.retroarch.rules
# This file is automatically generated.
# alsa
KERNEL=="controlC[0-9]*", TAG+="snap_retroarch_retroarch"
# alsa
KERNEL=="hwC[0-9]*D[0-9]*", TAG+="snap_retroarch_retroarch"
# alsa
KERNEL=="midiC[0-9]*D[0-9]*", TAG+="snap_retroarch_retroarch"
# alsa
KERNEL=="pcmC[0-9]*D[0-9]*[cp]", TAG+="snap_retroarch_retroarch"
# alsa
KERNEL=="seq", TAG+="snap_retroarch_retroarch"
# alsa
KERNEL=="timer", TAG+="snap_retroarch_retroarch"
# alsa
SUBSYSTEM=="sound", KERNEL=="card[0-9]*", TAG+="snap_retroarch_retroarch"
# bluez
KERNEL=="rfkill", TAG+="snap_retroarch_retroarch"
# joystick
KERNEL=="event[0-9]*", SUBSYSTEM=="input", ENV{ID_INPUT_JOYSTICK}=="1", TAG+="snap_retroarch_retroarch"
# joystick
KERNEL=="full", SUBSYSTEM=="mem", TAG+="snap_retroarch_retroarch"
# joystick
KERNEL=="js[0-9]*", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="dma_buf_te", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="galcore", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="mali[0-9]*", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="nvhost-*", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="nvmap", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="pvr_sync", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="vchiq", TAG+="snap_retroarch_retroarch"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_retroarch_retroarch"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="linux,cma", TAG+="snap_retroarch_retroarch"
# opengl
SUBSYSTEM=="dma_heap", KERNEL=="system", TAG+="snap_retroarch_retroarch"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_retroarch_retroarch"
# raw-usb
SUBSYSTEM=="tty", ENV{ID_BUS}=="usb", TAG+="snap_retroarch_retroarch"
# raw-usb
SUBSYSTEM=="usb", TAG+="snap_retroarch_retroarch"
# raw-usb
SUBSYSTEM=="usbmisc", TAG+="snap_retroarch_retroarch"
TAG=="snap_retroarch_retroarch", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_retroarch_retroarch $devpath $major:$minor"
